# deploy/telegram-private-bot-polling.service
[Unit]
Description=Telegram Private Bot Polling
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/telegram-private-bot-webhook
EnvironmentFile=/opt/telegram-private-bot-webhook/.env
ExecStart=/opt/telegram-private-bot-webhook/.venv/bin/telegram-private-bot-webhook
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
